# skills/calendar-plan/launchd/com.user.calendar-plan.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job for the calendar-plan skill (Claude target).

  Install:
    1. Edit <REPLACE_WITH_ABSOLUTE_PATH_TO_SKILL_DIR> below.
    2. cp this file to ~/Library/LaunchAgents/com.user.calendar-plan.plist
    3. launchctl load ~/Library/LaunchAgents/com.user.calendar-plan.plist
    4. Test: launchctl start com.user.calendar-plan
    5. Check logs: tail -f <SKILL_DIR>/logs/launchd.out

  Unload:
    launchctl unload ~/Library/LaunchAgents/com.user.calendar-plan.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.calendar-plan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string><REPLACE_WITH_ABSOLUTE_PATH_TO_SKILL_DIR>/claude/calendar-plan.sh</string>
  </array>

  <!-- Fires every day at 22:00 in the user's local timezone (TZ env below). -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Ensure timezone is consistent with what calendar-plan.sh expects. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>America/New_York</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string><REPLACE_WITH_ABSOLUTE_PATH_TO_SKILL_DIR>/claude/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string><REPLACE_WITH_ABSOLUTE_PATH_TO_SKILL_DIR>/claude/logs/launchd.err</string>

  <!-- Re-run within 5 minutes if a fire was missed (mac asleep / off). -->
  <key>StartCalendarIntervalRestartOnFail</key>
  <false/>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
